# tech-digest.service
[Unit]
Description=Daily Tech Digest Agent
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/daily_tech_digest
ExecStart=/usr/bin/python3 tech_digest_agent.py --schedule
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# 从 .env 文件加载环境变量
EnvironmentFile=/home/ubuntu/daily_tech_digest/.env

[Install]
WantedBy=multi-user.target
